# com.lwil.WATC.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lwil.WATC</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/username/Desktop/Loudness/WATC.sh</string>
    </array>
<key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>01</integer>
            <key>Weekday</key>
            <integer>6</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>23</integer>
            <key>Minute</key>
            <integer>01</integer>
            <key>Weekday</key>
            <integer>7</integer>
        </dict>
    </array>
</dict>
</plist>